# opodsync/rootfs/etc/s6-overlay/s6-rc.d/logs/run
#!/command/with-contenv bashio
# shellcheck shell=bash
bashio::log.info "Redirecting logs..."

exec /usr/bin/logs.sh
