# opodsync/rootfs/etc/s6-overlay/s6-rc.d/config/run
#!/command/with-contenv bashio
# shellcheck shell=bash
bashio::log.info "Configuring oPodSync..."

exec /usr/bin/config.sh
