# opodsync/rootfs/etc/s6-overlay/s6-rc.d/logrotate/run
#!/command/with-contenv bashio
# shellcheck shell=bash
bashio::log.info "Starting log rotation..."

exec /usr/bin/logrotate.sh
