# opodsync/rootfs/etc/s6-overlay/s6-rc.d/clone/run
#!/command/with-contenv bashio
# shellcheck shell=bash
bashio::log.info "Cloning oPodSync..."

exec /usr/bin/clone.sh
