# firefly-iii-dev/rootfs/etc/s6-overlay/s6-rc.d/init-firefly-iii/up
/etc/s6-overlay/s6-rc.d/init-firefly-iii/run